# automysqlbackup_check/check_amb.cron
# Install this file as /etc/cron.d/check_amb.cron

# check automysqlbackup data files once a week, every sunday at 3:00
0 3 * * 0 root /usr/local/bin/automysqlbackup_check >/var/lib/misc/automysqlbackup_check
